# kooka/etc/systemd/system/homeassistant.service.d/override.conf
[Service]
Environment=PATH=/opt/github.com/stuart12/python-scripts:/bin:/usr/bin:/usr/local/bin
