# kooka/etc/systemd/system/mpd.service.d/override.conf
[Unit]
After=pulseaudio.service
Wants=pulseaudio.service
[Service]
User=mpd
Group=pulse-access
Environment=PULSE_SINK=Snapcast
