# home-automation/home-assistant/hpmini110/snapclient.service.d/override.conf

# /etc/systemd/system/snapclient.service.d/override.conf
[Unit]
After=pulseaudio@homeassistant.service pulseaudio.service
